# piper-server/piper-tts.service
# Piper TTS HTTP Server - Systemd User Service
#
# Installation:
#   mkdir -p ~/.config/systemd/user
#   cp piper-tts.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable piper-tts
#   systemctl --user start piper-tts
#
# View logs:
#   journalctl --user -u piper-tts -f
#
# Stop:
#   systemctl --user stop piper-tts

[Unit]
Description=Piper TTS HTTP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/7DTD-LLM-NPC-Mod/piper-server
ExecStart=%h/.local/share/pipx/venvs/piper-tts/bin/python %h/7DTD-LLM-NPC-Mod/piper-server/piper_server.py --port 5050
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
